# crontab
# Run on container start
@start /usr/bin/download.sh

# Run hourly
@hourly /usr/bin/download.sh
